# templates/default/autochef.service.erb
# Dropped off by chef - changes will be overwritten.
# See https://github.com/PronghornDigital/mtna-server-cookbook

[Unit]
Description=Update Chef RPM and then run Chef

[Service]
Type=simple
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/bin/bash /opt/PronghornDigital/mtna-server-cookbook/autochef.sh
