# templates/default/mtna.service.erb
# Dropped off by chef - changes will be overwritten.
# See https://github.com/PronghornDigital/mtna-server-cookbook

[Unit]
Description=mtna node.js process
After=network-online.target

[Service]
User=mtna
Group=mtna
Restart=always
Environment=NODE_ENV=live
Environment=ARCHIVE_DATA_ROOT=/var/archives
ExecStart=/usr/local/bin/node /opt/PronghornDigital/mtna/dist/api/app.js

[Install]
WantedBy=multi-user.target
